# deploy/traiding.service
# Systemd service for Trading Gunicorn
# Copy to: /etc/systemd/system/trading.service
# Enable: sudo systemctl daemon-reload && sudo systemctl enable trading
# Start:  sudo systemctl start trading

[Unit]
Description=Trading Gunicorn daemon
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/trading/backend
ExecStart=/var/www/trading/backend/venv/bin/gunicorn config.wsgi:application -c gunicorn_config.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
